# wwwroot/python/crontab/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.FonswE/crontab installed on Wed Jan 29 12:19:10 2020)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
* * * * * echo "Hello world" >> /var/log/cron.log 2>&1
* * * * * echo $(date) >> /var/log/cron.log 2>&1
* * * * * echo "1111" >> /var/log/cron.log 2>&1
* * * * * echo $(date +\%T) >> /var/log/cron.log 2>&1stretch
